# tools/lc_notify/conf/my_app.ini
[program:flask_web_app]
command=python3 app.py ; 运行程序的命令
directory=/root/work/leetcode/tools/lc_notify ; 命令执行的目录
user=root ; 进程执行的用户身份
stopsignal=INT
autostart=true
autorestart=true ; 程序意外退出是否自动重启
startsecs=3 ; /自动重启间隔
stderr_logfile=/root/work/leetcode/tools/lc_notify/log/web_server_err.log ; 错误日志文件
stdout_logfile=/root/work/leetcode/tools/lc_notify/log/web_server_access.log ; 输出日志文件
